# docker/systemd/darioos-backup.timer
[Unit]
Description=Dario OS — agenda do backup diário (03:00)

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
