# docker/systemd/darioos-backup.service
[Unit]
Description=Dario OS — backup diário (PostgreSQL + Qdrant)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/projects/dario-os
ExecStart=%h/projects/dario-os/scripts/backup.sh
